# supervisord.conf
[supervisord]
nodaemon=true

[program:go-server]
command=go run main.go
directory=/app

[program:telethon-bot]
command=python3 bot.py
directory=/app
autostart=true
autorestart=true
